# .dotfiles/systemd/obsidian-headless-openboxsession.service.tmpl
[Unit]
Description=Headless openbox-session for Obsidian GUI (replaced by generic wm unit)
ConditionPathExists=false
