# .config/systemd/system/auto-update.service
[Unit]
Description=Automatic system update and reboot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto-update.sh
User=root
